# thelio/services.d/cron/run
#!/usr/bin/with-contenv bash

touch /var/log/cron.log
cron -f

